# deploy/crontab
10 * * * * BASH_ENV=$HOME/app/cron/env bash $HOME/app/cron/hourly.sh >> $HOME/cron.log 2>&1
20 0 * * * BASH_ENV=$HOME/app/cron/env bash $HOME/app/cron/daily.sh >> $HOME/cron.log 2>&1
